# cryptstore.service
[Unit]
Description=cryptstore deamon

[Service]
ExecStart=/usr/sbin/cryptstore

